# ansible/roles/jenkins/templates/jenkins.service.template
[Unit]
Description=Jenkins Daemon

[Service]
ExecStart=/usr/bin/java -jar /usr/local/etc/jenkins.war

[Install]
WantedBy=multi-user.target
